# chrome-extension/com.web_to_figma.ds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.web_to_figma.ds</string>
	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON_PATH__</string>
		<string>__INSTALL_DIR__/ds-daemon.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__INSTALL_DIR__</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__INSTALL_DIR__/daemon-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>__INSTALL_DIR__/daemon-stderr.log</string>
</dict>
</plist>
